# ops/systemd/quzijie-backup.timer
[Unit]
Description=Run Quzijie PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=15m
Unit=quzijie-backup.service

[Install]
WantedBy=timers.target
